# templates/disable_thp.service.j2
# /etc/systemd/system/disable_thp.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/disable_thp.sh

[Install]
WantedBy=multi-user.target
